# services/channel-monitor.service
[Unit]
Description=Claude Channel Monitor (Discord)
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/claude-autonomy-platform
Environment=PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/%i/claude-autonomy-platform/config/claude_infrastructure_config.txt
ExecStart=/usr/bin/python3 /home/%i/claude-autonomy-platform/discord/channel_monitor_simple.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target