# services/claude-auto-resume.service
[Unit]
Description=Claude Auto-Resume after System Reboot
After=network.target default.target
Documentation=file:///home/%i/claude-autonomy-platform/docs/AUTO_RESUME.md

[Service]
Type=oneshot
WorkingDirectory=/home/%i/claude-autonomy-platform
Environment=PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/%i/claude-autonomy-platform
EnvironmentFile=-/home/%i/claude-autonomy-platform/config/claude_infrastructure_config.txt
EnvironmentFile=-/home/%i/claude-autonomy-platform/config/claude.env
ExecStart=/home/%i/claude-autonomy-platform/core/claude_auto_resume.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=default.target
