# services/discord-status-bot.service
[Unit]
Description=Discord Status Bot for Claude
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/claude-autonomy-platform
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=%h/claude-autonomy-platform
EnvironmentFile=-%h/claude-autonomy-platform/config/claude_infrastructure_config.txt
ExecStart=/usr/bin/python3 %h/claude-autonomy-platform/discord/claude_status_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=process

[Install]
WantedBy=default.target