# services/discord-transcript-fetcher.service
[Unit]
Description=Discord Transcript Fetcher for ClAP
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/claude-autonomy-platform
ExecStart=/usr/bin/python3 %h/claude-autonomy-platform/discord/discord_transcript_fetcher.py
Restart=always
RestartSec=10

# Logging
StandardOutput=append:%h/claude-autonomy-platform/logs/transcript_fetcher.log
StandardError=append:%h/claude-autonomy-platform/logs/transcript_fetcher.log

[Install]
WantedBy=default.target
