# services/session-swap-monitor.service
[Unit]
Description=Session Swap Monitor for Claude
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/%i/claude-autonomy-platform
Environment=PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/%i/claude-autonomy-platform
EnvironmentFile=-/home/%i/claude-autonomy-platform/config/claude_infrastructure_config.txt
ExecStart=/usr/bin/python3 /home/%i/claude-autonomy-platform/core/session_swap_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=process
EnvironmentFile=-/home/%i/claude-autonomy-platform/config/claude.env

[Install]
WantedBy=default.target