# services/session-bridge-monitor.service
[Unit]
Description=Session Bridge Monitor for Claude
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/%i/claude-autonomy-platform
ExecStart=/usr/bin/python3 /home/%i/claude-autonomy-platform/session_bridge_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/home/%i/claude-autonomy-platform/config/claude.env

[Install]
WantedBy=default.target