# systemd/led-daemon.service
[Unit]
Description=LED State Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/claude-autonomy-platform
Environment="CLAP_DIR=%h/claude-autonomy-platform"
Environment="PYTHONUNBUFFERED=1"
ExecStart=%h/claude-autonomy-platform/core/led_daemon.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
